# etc/init/mountain-tai.conf
description "Mountain Tai"
author "frazy <frazy@thstack.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start --exec /usr/sbin/mountain-scheduler
